# flask-omr.service
[Unit]
Description=Flask OMR AI Reader Service
After=network.target

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/var/www/html/florescer-ia
Environment="PATH=/root/.cargo/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/.cargo/bin/uv run gunicorn --bind 127.0.0.1:5000 --workers 4 --timeout 300 --access-logfile /var/log/flask-omr-access.log --error-logfile /var/log/flask-omr-error.log app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

